# containers/openconnect/rootfs/etc/s6-overlay/s6-rc.d/sing-box/run
#!/command/execlineb -P
sing-box -D /var/lib/sing-box -C /etc/sing-box run
